# supervisord.conf
[supervisord]
nodaemon = true

[program: cloudflared]
command=/opt/argo/cloudflared
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
